# eraseScreen.service
[Unit]
Description=Bitcoin daemon

[Service]
ExecStart=/home/pi/iHodl/eraseScreen.sh

[Install]
WantedBy=multi-user.target

